# supervisor-vault.conf
[program:vault]
command=/app/code/vault server -config=/app/data/config.hcl
user=cloudron
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
